# crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# Cron файл для ежедневной суммаризации
# Формат: минута час день месяц день_недели команда
# Запускаем скрипт каждые 5 минут, чтобы он сам проверял настройки времени для каждого чата
# из config.yaml (summarize.daily_time). Это позволяет разным чатам иметь разное время запуска.
# Запуск проверки ежедневной суммаризации каждые 5 минут
# Скрипт сам проверит настройки времени для каждого чата из config.yaml
*/5 * * * * /app/run_daily_summary.sh
